# chapter7/section3/ansible/dev/redis/templates/upstart.conf
description "redis server"

start on runlevel [23]
stop on shutdown

exec sudo -u redis /usr/local/bin/redis-server /etc/redis/redis.conf

respawn
respawn limit 5 60